# systemd/user/emacs.socket
[Socket]
ListenStream=~/.emacs.socket
DirectoryMode=0700

[Install]
WantedBy=sockets.target
